# autostart/walsh.desktop
[Desktop Entry]
Name=Wallpaper Shuffler
Type=Application
Exec=walsh
TryExec=walsh
StartupNotify=false
